# quaffle/quaffle-ws-service/etc/service/quaffle-ws-service/run
#!/bin/bash -e

PACKAGE="quaffle-ws-service"

CP="/usr/share/$PACKAGE/$PACKAGE.jar"

java -jar $CP 2>&1
