# quaffle/quaffle-ws-service/DEBIAN/postinst
#!/bin/bash -e

PAC="quaffle-ws"

invoke-rc.d $PAC start
