# quaffle-frontend/quaffle-frontend-service/etc/service/quaffle-frontend-service/run
#!/bin/bash -e

PACKAGE="quaffle-frontend-service"

CP="/usr/share/$PACKAGE/$PACKAGE.js"

exec node $CP 2>&1
